# rules/sigma/mitre/ta0003_persistence/T1543_003_windows_service.yml
title: Windows Service Created for Persistence
id: c3d4e5f6-a7b8-9012-cdef-123456789012
status: stable
description: |
  Detects creation of new Windows services, which is a common persistence
  mechanism used by malware and attackers after gaining initial access.
  Covers EventID 7045 (Service Control Manager) and schtasks abuse.
references:
  - https://attack.mitre.org/techniques/T1543/003/
  - https://docs.microsoft.com/en-us/windows/security/threat-protection/auditing/event-7045
author: Chronos-DFIR / Ivan Huerta
date: 2026/03/03
tags:
  - attack.persistence
  - attack.t1543.003
  - attack.privilege_escalation
logsource:
  category: service_install
  product: windows
detection:
  new_service:
    EventID: 7045
  suspicious_service_name:
    ServiceName|contains|any:
      - 'update'
      - 'svc'
      - 'host'
      - 'agent'
      - 'backup'
  suspicious_binary_path:
    ImagePath|contains|any:
      - '\Temp\'
      - '\AppData\'
      - '\Users\Public\'
      - 'powershell'
      - 'cmd.exe /c'
      - 'wscript'
      - 'cscript'
      - 'mshta'
  condition: new_service and suspicious_binary_path
falsepositives:
  - Legitimate software installers creating services
  - Known backup or monitoring agents
level: medium
fields:
  - ServiceName
  - ImagePath
  - ServiceType
  - StartType
custom:
  mitre_tactic: "TA0003 – Persistence"
  mitre_technique: "T1543.003 – Windows Service"
  windows_event_id: "7045"
